# deploy/dashboard.service
[Unit]
Description=Home Dashboard eInk renderer
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=oneshot
# These values are filled in automatically by `make pi-enable`.
# __USER__ and __INSTALL_DIR__ are substituted with the current user and repo path.
User=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=GPIOZERO_PIN_FACTORY=lgpio
ExecStart=__INSTALL_DIR__/venv/bin/python -m src.main --config config/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=append:__INSTALL_DIR__/output/dashboard.log
StandardError=append:__INSTALL_DIR__/output/dashboard.log
